# modules/AlarmClock/alarms.cron
37	20	*	*	*	/home/mde/bin/gladQuotes.sh
